# services/ogar/run
#!/bin/bash
cd /app/ogar
node /app/ogar/src/index.js
